# tungstite.service
[Unit]
Description=Tungstite IRC Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=reverse
WorkingDirectory=/home/reverse/irc/scripts/tungstite
ExecStart=/home/reverse/irc/scripts/tungstite/.env/bin/python -m tungstite config.yaml
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
